# ubuntu/etc/profile.d/aliases.sh
alias cp='cp -i'
alias mv='mv -i'
alias rm='rm -i'
alias e='subl'
alias flushdns='sudo service nscd restart'
alias storm='/usr/local/phpstorm/bin/phpstorm.sh'
