# ubuntu/etc/profile.d/dragondisk.sh
alias dragondisk='cd /usr/local/dragondisk/ ; ./dragondisk'
